# dynamicbg.service
[Unit]
Description=dynamic-bg

[Service]
ExecStart=/usr/bin/ArcgDynamicBackground
[Install]
WantedBy=multi-user.target 
